# crontab
*/5 * * * * /root/rft.sh >/tmp/rft.sh.out
